# scripts/95-monitor-hotplug.rules
# =============================================================================
# 95-monitor-hotplug.rules
# Regla udev para detección automática de monitor externo (VGA/HDMI/DP)
# =============================================================================
#
# INSTALACIÓN:
#   sudo cp 95-monitor-hotplug.rules /etc/udev/rules.d/
#   sudo udevadm control --reload-rules
#
# IMPORTANTE: Edita los dos campos antes de copiar:
#   - NOMBRE_USUARIO → tu usuario real (ej: "mario")
#   - /home/NOMBRE_USUARIO → ruta real a tu home
#
# =============================================================================

ACTION=="change", SUBSYSTEM=="drm", \
  RUN+="/bin/su mat -c 'DISPLAY=:0 XAUTHORITY=/home/mat/.Xauthority /home/mat/.config/i3/scripts/display.sh'"